# contentops/coverage/data/mitre_attack_techniques.json
{
  "schema_version": "1.0",
  "source": "MITRE ATT&CK Enterprise — curated subset",
  "note": "Curated reference list of high-value techniques across all 14 tactics. NOT exhaustive — supply --techniques-file with a fuller MITRE Enterprise STIX-derived list for full coverage analysis. Tactics use the pipeline's PascalCase convention to match metadata.tactics in detection envelopes.",
  "techniques": [
    {"id": "T1595", "name": "Active Scanning", "tactics": ["Reconnaissance"]},
    {"id": "T1592", "name": "Gather Victim Host Information", "tactics": ["Reconnaissance"]},
    {"id": "T1589", "name": "Gather Victim Identity Information", "tactics": ["Reconnaissance"]},
    {"id": "T1590", "name": "Gather Victim Network Information", "tactics": ["Reconnaissance"]},
    {"id": "T1591", "name": "Gather Victim Org Information", "tactics": ["Reconnaissance"]},
    {"id": "T1583", "name": "Acquire Infrastructure", "tactics": ["ResourceDevelopment"]},
    {"id": "T1586", "name": "Compromise Accounts", "tactics": ["ResourceDevelopment"]},
    {"id": "T1587", "name": "Develop Capabilities", "tactics": ["ResourceDevelopment"]},
    {"id": "T1588", "name": "Obtain Capabilities", "tactics": ["ResourceDevelopment"]},
    {"id": "T1078", "name": "Valid Accounts", "tactics": ["DefenseEvasion", "InitialAccess", "Persistence", "PrivilegeEscalation"]},
    {"id": "T1190", "name": "Exploit Public-Facing Application", "tactics": ["InitialAccess"]},
    {"id": "T1133", "name": "External Remote Services", "tactics": ["InitialAccess", "Persistence"]},
    {"id": "T1566", "name": "Phishing", "tactics": ["InitialAccess"]},
    {"id": "T1199", "name": "Trusted Relationship", "tactics": ["InitialAccess"]},
    {"id": "T1195", "name": "Supply Chain Compromise", "tactics": ["InitialAccess"]},
    {"id": "T1059", "name": "Command and Scripting Interpreter", "tactics": ["Execution"]},
    {"id": "T1053", "name": "Scheduled Task/Job", "tactics": ["Execution", "Persistence", "PrivilegeEscalation"]},
    {"id": "T1106", "name": "Native API", "tactics": ["Execution"]},
    {"id": "T1129", "name": "Shared Modules", "tactics": ["Execution"]},
    {"id": "T1204", "name": "User Execution", "tactics": ["Execution"]},
    {"id": "T1098", "name": "Account Manipulation", "tactics": ["Persistence", "PrivilegeEscalation"]},
    {"id": "T1136", "name": "Create Account", "tactics": ["Persistence"]},
    {"id": "T1547", "name": "Boot or Logon Autostart Execution", "tactics": ["Persistence", "PrivilegeEscalation"]},
    {"id": "T1543", "name": "Create or Modify System Process", "tactics": ["Persistence", "PrivilegeEscalation"]},
    {"id": "T1546", "name": "Event Triggered Execution", "tactics": ["Persistence", "PrivilegeEscalation"]},
    {"id": "T1505", "name": "Server Software Component", "tactics": ["Persistence"]},
    {"id": "T1068", "name": "Exploitation for Privilege Escalation", "tactics": ["PrivilegeEscalation"]},
    {"id": "T1134", "name": "Access Token Manipulation", "tactics": ["DefenseEvasion", "PrivilegeEscalation"]},
    {"id": "T1484", "name": "Domain Policy Modification", "tactics": ["DefenseEvasion", "PrivilegeEscalation"]},
    {"id": "T1055", "name": "Process Injection", "tactics": ["DefenseEvasion", "PrivilegeEscalation"]},
    {"id": "T1027", "name": "Obfuscated Files or Information", "tactics": ["DefenseEvasion"]},
    {"id": "T1070", "name": "Indicator Removal", "tactics": ["DefenseEvasion"]},
    {"id": "T1218", "name": "System Binary Proxy Execution", "tactics": ["DefenseEvasion"]},
    {"id": "T1562", "name": "Impair Defenses", "tactics": ["DefenseEvasion"]},
    {"id": "T1497", "name": "Virtualization/Sandbox Evasion", "tactics": ["DefenseEvasion", "Discovery"]},
    {"id": "T1112", "name": "Modify Registry", "tactics": ["DefenseEvasion"]},
    {"id": "T1140", "name": "Deobfuscate/Decode Files or Information", "tactics": ["DefenseEvasion"]},
    {"id": "T1003", "name": "OS Credential Dumping", "tactics": ["CredentialAccess"]},
    {"id": "T1110", "name": "Brute Force", "tactics": ["CredentialAccess"]},
    {"id": "T1555", "name": "Credentials from Password Stores", "tactics": ["CredentialAccess"]},
    {"id": "T1056", "name": "Input Capture", "tactics": ["CredentialAccess", "Collection"]},
    {"id": "T1558", "name": "Steal or Forge Kerberos Tickets", "tactics": ["CredentialAccess"]},
    {"id": "T1552", "name": "Unsecured Credentials", "tactics": ["CredentialAccess"]},
    {"id": "T1018", "name": "Remote System Discovery", "tactics": ["Discovery"]},
    {"id": "T1083", "name": "File and Directory Discovery", "tactics": ["Discovery"]},
    {"id": "T1057", "name": "Process Discovery", "tactics": ["Discovery"]},
    {"id": "T1087", "name": "Account Discovery", "tactics": ["Discovery"]},
    {"id": "T1082", "name": "System Information Discovery", "tactics": ["Discovery"]},
    {"id": "T1016", "name": "System Network Configuration Discovery", "tactics": ["Discovery"]},
    {"id": "T1069", "name": "Permission Groups Discovery", "tactics": ["Discovery"]},
    {"id": "T1021", "name": "Remote Services", "tactics": ["LateralMovement"]},
    {"id": "T1570", "name": "Lateral Tool Transfer", "tactics": ["LateralMovement"]},
    {"id": "T1080", "name": "Taint Shared Content", "tactics": ["LateralMovement"]},
    {"id": "T1550", "name": "Use Alternate Authentication Material", "tactics": ["DefenseEvasion", "LateralMovement"]},
    {"id": "T1005", "name": "Data from Local System", "tactics": ["Collection"]},
    {"id": "T1114", "name": "Email Collection", "tactics": ["Collection"]},
    {"id": "T1213", "name": "Data from Information Repositories", "tactics": ["Collection"]},
    {"id": "T1115", "name": "Clipboard Data", "tactics": ["Collection"]},
    {"id": "T1071", "name": "Application Layer Protocol", "tactics": ["CommandAndControl"]},
    {"id": "T1573", "name": "Encrypted Channel", "tactics": ["CommandAndControl"]},
    {"id": "T1095", "name": "Non-Application Layer Protocol", "tactics": ["CommandAndControl"]},
    {"id": "T1105", "name": "Ingress Tool Transfer", "tactics": ["CommandAndControl"]},
    {"id": "T1219", "name": "Remote Access Software", "tactics": ["CommandAndControl"]},
    {"id": "T1090", "name": "Proxy", "tactics": ["CommandAndControl"]},
    {"id": "T1568", "name": "Dynamic Resolution", "tactics": ["CommandAndControl"]},
    {"id": "T1041", "name": "Exfiltration Over C2 Channel", "tactics": ["Exfiltration"]},
    {"id": "T1567", "name": "Exfiltration Over Web Service", "tactics": ["Exfiltration"]},
    {"id": "T1029", "name": "Scheduled Transfer", "tactics": ["Exfiltration"]},
    {"id": "T1486", "name": "Data Encrypted for Impact", "tactics": ["Impact"]},
    {"id": "T1490", "name": "Inhibit System Recovery", "tactics": ["Impact"]},
    {"id": "T1485", "name": "Data Destruction", "tactics": ["Impact"]},
    {"id": "T1499", "name": "Endpoint Denial of Service", "tactics": ["Impact"]},
    {"id": "T1498", "name": "Network Denial of Service", "tactics": ["Impact"]},
    {"id": "T1489", "name": "Service Stop", "tactics": ["Impact"]},
    {"id": "T1531", "name": "Account Access Removal", "tactics": ["Impact"]}
  ]
}
